# contrib/init/clashicd.conf
description "Title Network Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env TITLED_BIN="/usr/bin/titled"
env TITLED_USER="bitcoin"
env TITLED_GROUP="bitcoin"
env TITLED_PIDDIR="/var/run/titled"
# upstart can't handle variables constructed with other variables
env TITLED_PIDFILE="/var/run/titled/titled.pid"
env TITLED_CONFIGFILE="/etc/title/title.conf"
env TITLED_DATADIR="/var/lib/titled"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # titled will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$TITLED_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run titled."
        echo "The setting must appear in $TITLED_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"TItle Network Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$TITLED_PIDDIR"
    chmod 0755 "$TITLED_PIDDIR"
    chown $TITLED_USER:$TITLED_GROUP "$TITLED_PIDDIR"
    chown $TITLED_USER:$TITLED_GROUP "$TITLED_CONFIGFILE"
    chmod 0660 "$TITLED_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$TITLED_PIDFILE" \
    --chuid $TITLED_USER:$TITLED_GROUP \
    --exec "$TITLED_BIN" \
    -- \
    -pid="$TITLED_PIDFILE" \
    -conf="$TITLED_CONFIGFILE" \
    -datadir="$TITLED_DATADIR" \
    -disablewallet \
    -daemon

